# environment/sber-together-api.crontab
PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin
*/30 * * * * /opt/sber-together-api/current/app/scripts/cron/checkOrderStatus.sh
0 12 * * * /opt/sber-together-api/current/app/scripts/cron/monthlyPayments.sh
